# meta-oniro-blueprints/recipes-apps/vending-machine/vending-machine-control-application/vending-machine-control-application.service
# SPDX-FileCopyrightText: Huawei Inc.
#
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Vending Machine Control Application

[Service]
Type=notify
Environment="VENDING_MACHINE_CONTROL_I2C_BUS=@VENDING_MACHINE_CONTROL_I2C_BUS@"
ExecStart=/usr/bin/vending-machine-control-application
Restart=always

[Install]
WantedBy=multi-user.target
