# meta-oniro-staging/recipes-core/base-files/base-files/profile.d/tmout.sh
# Shell timeout definition

TMOUT=600
readonly TMOUT
export TMOUT
